# python/miniconda3_py313/etc/profile.d/conda.csh
setenv CONDA_EXE "/home/rreis/python/miniconda3_py313/bin/conda";
setenv _CONDA_EXE "/home/rreis/python/miniconda3_py313/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/home/rreis/python/miniconda3_py313/bin/python";
setenv _CONDA_ROOT "/home/rreis/python/miniconda3_py313";

source "/home/rreis/python/miniconda3_py313/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"